# crontab
# RSI - Crypto (todos os dias às 21:00 BR)
00 00 * * * /app/rsi_screening.sh crypto

# RSI - Ações/ETF (seg-sex às 20:30 BR)
00 23 * * 1-5 /app/rsi_screening.sh acoes
